# vagrant/ansible/roles/common/files/etc/profile.d/PS1.sh
PS1="\[\033[36m\][\[\033[m\]\u@\H \W\[\033[36m\]]\\$\[\033[m\] "
